# scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
command=/startup.sh